# systemd/servicematrix-frontend.service
[Unit]
Description=ServiceMatrix Frontend (Next.js)
After=network.target servicematrix-backend.service
Wants=servicematrix-backend.service

[Service]
Type=simple
User=servicematrix
WorkingDirectory=/opt/servicematrix/frontend
EnvironmentFile=/opt/servicematrix/.env
ExecStartPre=/opt/servicematrix/scripts/linux/check-ports.sh
ExecStart=/opt/servicematrix/scripts/linux/start-frontend.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=servicematrix-frontend

[Install]
WantedBy=multi-user.target
